# misc/etc/cron.d/logout_phones
# Logout tagged phones
23	1	*	*	*	root	/opt/GS5/script/logout_phones.sh

